# loki.conf
[program:loki]
command=/usr/bin/loki --config.file=/etc/loki/loki.yaml
priority=1
redirect_stderr=true
