# vector.conf
[program:vector]
command=/usr/bin/vector -c /etc/vector/vector.toml
priority=2
redirect_stderr=true
